# root/etc/services.d/oidctestserver/run
#!/bin/execlineb -P
#cd /opt/oidc/oidctestserver/oidc_op
# sh run.sh
